# crontab
*/1 * * * * python -u /app/update-home-dns.py
