# crontab
*/1 * * * * python -u /update-home-dns.py
