# crontab
0 * * * * python -u /app/update-home-dns.py
